# deploy/drynn.service
[Unit]
Description=Hobo web application
Documentation=https://github.com/mdhender/drynn
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=drynn
Group=drynn

WorkingDirectory=/opt/drynn
ExecStart=/opt/drynn/bin/server --config /var/lib/drynn/config/server.json
EnvironmentFile=-/etc/drynn/drynn.env

Restart=on-failure
RestartSec=5s
TimeoutStopSec=15s

# Logging: use the journal.
StandardOutput=journal
StandardError=journal
SyslogIdentifier=drynn

# Sandboxing.
NoNewPrivileges=true
PrivateTmp=true
PrivateDevices=true
ProtectSystem=strict
ProtectHome=true
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectKernelLogs=true
ProtectControlGroups=true
ProtectClock=true
ProtectHostname=true
ProtectProc=invisible
RestrictNamespaces=true
RestrictRealtime=true
RestrictSUIDSGID=true
LockPersonality=true
MemoryDenyWriteExecute=true
SystemCallArchitectures=native
SystemCallFilter=@system-service
SystemCallFilter=~@privileged @resources

# /var/lib/drynn must be writable; everything else is read-only.
ReadWritePaths=/var/lib/drynn

CapabilityBoundingSet=
AmbientCapabilities=

# Resource limits.
LimitNOFILE=65536

[Install]
WantedBy=multi-user.target
